# deploy/demo-app.path
[Path]
PathModified=/opt/deployments/restart.flag
Unit=demo-app-restart.service

[Install]
WantedBy=multi-user.target
